# Gogios/DEBIAN/prerm
#!/bin/bash
# Author: Bailey Kasin

systemctl stop service-check-web
systemctl stop service-check
systemctl disable service-check-web
systemctl disable service-check