# Gogios/DEBIAN/prerm
#!/bin/bash
# Author: Bailey Kasin

systemctl stop service-check
systemctl disable service-check

userdel gogios