# Gogios/etc/systemd/system/gogios-web.service
[Unit]
Description=Controls service check engine website
After=network.target
[Service]
Type=simple
User=root
ExecStart=dotnet run --project /opt/site
[Install]
WantedBy=multi-user.target